# openclaw-adapter/systemd/wf-openclaw-adapter.service
[Unit]
Description=WF OpenClaw Adapter Service
After=network.target
Wants=network.target

[Service]
# 核心启动命令（需替换为实际的 jar 包绝对路径）
ExecStart=/usr/bin/java -jar /opt/openclaw/openclaw-adapter-1.0.0.jar
# 可选：添加 JVM 内存参数（根据需求启用）
# ExecStart=/usr/bin/java -Xms512m -Xmx1024m -jar /opt/openclaw/openclaw-adapter-1.0.0.jar

# 优雅停止服务
ExecStop=/bin/kill -15 $MAINPID
# 进程异常退出后自动重启
Restart=on-failure
RestartSec=5s

# 日志输出配置
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target